# rootfs/ap121-2.6.31-2MB/etc/inittab
::sysinit:/etc/rc.d/rcS
::respawn:/sbin/getty ttyS0 115200
::shutdown:/bin/cfg wlan down
