# rootfs/ap91-2MB/etc/inittab
::sysinit:/etc/rc.d/rcS
::respawn:/bin/sh -
::shutdown:/bin/cfg wlan down
